# iot-certbot/rootfs/etc/services.d/lets-encrypt/finish
#!/usr/bin/env bashio
# -*- shell-script -*-
# ==============================================================================
# Take down the S6 supervision tree when letsencrypt fails
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.warning "Halt add-on"
/run/s6/basedir/bin/halt
