# iot-certbot/rootfs/etc/services.d/lets-encrypt/run
#!/usr/bin/with-contenv bashio
# -*- shell-script -*-
# ==============================================================================
# Re-run the renew script forever, until the addon is stopped. Hourly ought to
# be good enough.
# ==============================================================================
while true; do
  /opt/letsencrypt/bin/renew || bashio::log.error "Most recent letsencrypt run failed; check your logs"
  sleep 3600
done
