# rootfs/etc/periodic/15min/update_iot_creds
#!/usr/bin/env bash
/opt/iot/bin/update-credentials.sh
